# etc/systemd/system/disable-wakeup.service
[Unit]
Description=Disable ACPI wake device

[Service]
Type=oneshot
ExecStart=/usr/local/bin/disable-wakeup.sh

[Install]
WantedBy=multi-user.target
